# ssh_updater_v24_10_08_2017/rc.local
ethtool -s eth0 duplex half speed 10 autoneg off
